# etc/udev/rules.d/81-wifi-powersave.rules
SUBSYSTEM=="net", KERNEL=="wlan[0-9]", ACTION=="add", RUN+="/usr/bin/iw dev $name set power_save off"
